# gradle/dist/kube-podman.service
[Unit]
Description=Podman API Service
Requires=kube-podman.socket
After=kube-podman.socket
Documentation=man:podman-system-service(1)
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/podman system service
KillMode=process

[Install]
WantedBy=multi-user.target
Also=kube-podman.socket
